# scripts/rc.local
/usr/sbin/daemon -rft "IPV4 Keepalive" /sbin/ping 8.8.8.8

